# ansible-geonode/templates/supervisor_flower.conf.j2
[program:flower]
command={{ celery_scripts_dir }}/flower_start

autostart=true
autorestart=true

user={{ celery_user }}